# infra/ansible/roles/kafka_kraft/templates/kafka.service.j2
# -----------------------------------------------------------------------------
# file # /usr/lib/systemd/system/kafka.service for {{ ansible_hostname }} ({{ ansible_host }})
# {{ ansible_managed }}
# generated by {{ role_name }} role from {{ template_path | basename }}
# at {{ template_run_date }} on {{ template_host }}
# -----------------------------------------------------------------------------
#

[Unit]
Description=Apache Kafka Kraft
Documentation=http://kafka.apache.org/documentation.html
# After=network-online.target
# Requires=network-online.target
After=network-online.target remote-fs.target
Requires=network-online.target remote-fs.target

[Service]
Type=simple
## Restart=on-failure
User=kafka
Group=kafka
### Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk-amd64
Environment=JAVA_HOME=/usr/lib/jvm/java-21-openjdk-amd64
### ExecStart={{ kafka_home }}/kafka/bin/kafka-server-start.sh {{ kafka_home }}/kafka/config/kraft/server.properties
ExecStart={{ kafka_home }}/kafka/bin/kafka-server-start.sh {{ kafka_home }}/kafka/config/server.properties
ExecStop={{ kafka_home }}/kafka/bin/kafka-server-stop.sh
WorkingDirectory={{ kafka_home }}

[Install]
WantedBy=multi-user.target

##

